# xmg-kb.service
[Unit]
Description=XMG Keyboard RGB Controller
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xmg-kb --restore
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

